# flask
#!/bin/bash
# chkconfig: - 99 10
python /home/ec2-user/main.py > /tmp/main.log 2>&1 &
